# package/udev/rules.d/10-wlan-restore.rules
KERNEL=="wlan0", RUN+="/etc/virt2real/wlanrestore.sh"

